# daemon/hddwatch/DEBIAN/prerm
#!/bin/sh

systemctl disable hddwatch.timer

exit 0