# daemons/hddwatch/DEBIAN/postinst
#!/bin/sh

systemctl enable hddwatch.timer

exit 0
