# daemon/fancontrol/fancontrol.service
[Unit]
Description=Fan controller 
OnFailure=device_failure.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/fancontrol
