# daemons/hddwatch/lib/systemd/system/hddwatch.timer
[Unit]
Description=Run hddmonitor every 6 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=5min
Unit=hddwatch.service

[Install]
WantedBy=timers.target