# daemons/hddwatch/lib/systemd/system/hddwatch.service


[Unit]
Description=Watch HDDs for reallocated sectors 

[Service]
ExecStart=/usr/bin/hddwatch